# etc/configuration-server.service
[Unit]
Description=configuration-server
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
User=www-data
Group=www-data
WorkingDirectory=/opt/configuration-server
ExecStart=bash /opt/configuration-server/start.sh production

[Install]
WantedBy=default.target
